# deploy/feedling-mcp.service
[Unit]
Description=Feedling MCP Server
After=network.target feedling-backend.service
Requires=feedling-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/feedling-mcp-v1/backend
EnvironmentFile=/home/ubuntu/feedling.env
ExecStart=/home/ubuntu/feedling-venv/bin/python mcp_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
